# images/browser/files/etc/systemd/system/qemu-guest-agent.service.d/aegis-lockdown.conf
# Aegis — qemu-guest-agent lockdown (spec §4: brak guest tools z integracjami).
#
# Identical policy to the gateway: the agent is kept ONLY for clean ACPI shutdown
# and fsfreeze/thaw (so overlay teardown is consistent). Every command that could
# leak host data or provide desktop/file integration is disabled via the RPC
# allow-list. There is no clipboard, no file transfer, no guest-exec.
[Service]
ExecStart=
ExecStart=/usr/sbin/qemu-ga \
  --allow-rpcs=guest-ping,guest-sync,guest-sync-delimited,guest-shutdown,guest-fsfreeze-freeze,guest-fsfreeze-thaw,guest-fsfreeze-status
NoNewPrivileges=yes
ProtectHome=yes
