# images/browser/files/etc/systemd/system/aegis-kiosk.service.d/10-hardening.conf
# Aegis Browser — seccomp / systemd sandbox hardening drop-in for the kiosk
# (spec §10: "filtrowanie wywołań systemowych", "sandbox renderera pozostaje
# aktywny").
#
# IMPORTANT: Chromium ships its OWN seccomp-bpf sandbox for its renderer/GPU/
# utility processes, layered on top of setuid + user-namespace sandboxing. That
# inner sandbox is the primary defense and MUST NOT be disabled (never
# --no-sandbox). This drop-in adds a systemd-level seccomp allow-list at the
# SERVICE boundary — it is intentionally permissive enough to let Chromium
# install its own, stricter filters, while still blocking whole classes of
# syscalls a compromised launcher would want (module loading, kexec, raw ptrace
# of other processes, obsolete/uncommon syscalls, clock/time manipulation).
#
# We use @system-service as the base (which permits the syscalls Chromium needs,
# including the seccomp() and clone/unshare calls its sandbox bootstrap uses) and
# then subtract dangerous groups. This keeps Chromium fully functional while
# removing host-attack surface if the renderer sandbox were ever escaped.
[Service]
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@module @reboot @swap @raw-io @obsolete @cpu-emulation
SystemCallFilter=~@clock @debug
SystemCallErrorNumber=EPERM
# Keep user namespaces available (Chromium sandbox needs them); do not add
# RestrictNamespaces=yes here. clone/unshare/setns remain permitted by
# @system-service.
