# images/gateway/files/etc/systemd/system/qemu-guest-agent.service.d/aegis-lockdown.conf
# Aegis — qemu-guest-agent lockdown (spec §4: brak guest tools z integracjami).
#
# qemu-guest-agent is retained ONLY so libvirt can request a clean ACPI shutdown
# and freeze/thaw for consistent overlay teardown. Every command that could leak
# host data or provide desktop/file integration is blocked here. The agent runs
# with an explicit allow-list; guest-exec, file read/write, and clipboard-adjacent
# commands are disabled.
[Service]
# Restrict the agent to a minimal command set (shutdown/ping/sync/fsfreeze only).
ExecStart=
ExecStart=/usr/sbin/qemu-ga \
  --allow-rpcs=guest-ping,guest-sync,guest-sync-delimited,guest-shutdown,guest-fsfreeze-freeze,guest-fsfreeze-thaw,guest-fsfreeze-status
NoNewPrivileges=yes
ProtectHome=yes
